# dev-python/graphite-api/files/init.d-graphite-api
#!/sbin/openrc-run

description="Graphite API Daemon"
supervisor=supervise-daemon
respawn_period="${respawn_period:-30}"
term_type="${term_type:-linux}"
command=/usr/bin/gunicorn
pidfile="/run/${RC_SVCNAME}.pid"
command_user="nobody:nobody"
command_args="graphite_api.app:app -w${GRAPHITE_API_WORKERS} -b ${GRAPHITE_API_PORT}"

depend() {
        need net
}

start_pre() {
        export EINFO_QUIET="${quiet:-yes}"
}

stop_pre()
{
        export EINFO_QUIET="${quiet:-yes}"
}
